# root/etc/services.d/calibre/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc /opt/calibre/calibre-server \
  --port 8080 \
  --enable-local-write \
  /books